# pi-setup/sangi-notification-monitor@.service
[Unit]
Description=SANGI Notification Monitor Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/SANGI/PC-setup/notification-service
ExecStart=/usr/bin/python3 /home/%i/SANGI/PC-setup/notification-service/notification_service.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment="PYTHONUNBUFFERED=1"
Environment="DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus"

[Install]
WantedBy=default.target
